# packages/promptflow-image-bundle/promptflow_image_bundle-0.3.0-py3-none-any.whl/promptflow/runtime/package_data/runit/embedding-store-service/finish
#!/bin/bash

echo "`date -uIns` - embedding-store-service/finish $@"
